# system-rootfs/deck/shared/usr/lib/systemd/system/critos-autologin.service
[Unit]
Description=Enables CritOS autologin
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/critos-autologin

[Install]
WantedBy=multi-user.target
